# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24
autorestart=true
priority=1
startretries=3
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

[program:chromium]
command=chromium --no-sandbox --disable-setuid-sandbox --disable-dev-shm-usage --disable-gpu --remote-debugging-port=9222 --remote-debugging-address=127.0.0.1 --user-data-dir=/data/chrome-profile --display=:99 --window-size=1920,1080 --start-maximized --no-first-run --disable-background-timer-throttling --disable-renderer-backgrounding --password-store=basic --use-mock-keychain
environment=DISPLAY=":99"
autorestart=true
priority=2
startretries=5
stdout_logfile=/var/log/supervisor/chromium.log
stderr_logfile=/var/log/supervisor/chromium.log

[program:x11vnc]
command=x11vnc -display :99 -localhost -nopw -forever -shared -rfbport 5900 -noxdamage -xkb
environment=DISPLAY=":99"
autorestart=true
priority=3
startretries=3
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

[program:novnc]
command=websockify --web=/usr/share/novnc --cert=none 127.0.0.1:8081 127.0.0.1:5900
autorestart=true
priority=4
startretries=3
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log

[program:openclaw-wrapper]
command=node src/server.js
directory=/app
# Don't auto-restart forever - let Railway restart the container on persistent failure
autorestart=unexpected
exitcodes=0,2
startretries=3
startsecs=5
priority=5
stdout_logfile=/var/log/supervisor/wrapper.log
stderr_logfile=/var/log/supervisor/wrapper.log
environment=BROWSER_CDP_URL="http://localhost:9222"
